# docs/supervisor/ramassage_admin.conf

[program:ramassage_admin]
#environment = PATH="/usr/local/python-3.4.3/bin/:%(ENV_PATH)s"
command = /usr/local/python-3.4.3/bin/gunicorn admin:app -b localhost:9000
directory = /home/ramassage/ramassage.epitech.eu/
user = ramassage
